# deploy/jeraaiboss/pulse-restore-test.service
[Unit]
Description=Pulse Postgres backup restore verification

[Service]
Type=oneshot
ExecStart=%h/pulse-db/pulse-restore-test.sh
