# deploy/jeraaiboss/pulse-backup.service
[Unit]
Description=Pulse Postgres encrypted backup
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/pulse-db/pulse-backup.sh
